# scripts/services/nodebb.conf
start on startup
stop on runlevel [016]
respawn
setuid nodebb
setgid nogroup
script
        cd /path/to/nodebb/
        node loader.js --no-silent --no-daemon
end script
