# scripts/services/notifico.conf
description "Notifico daemon"

start on started networking
stop on runlevel [!2345]

setuid notifico
setgid nogroup

chdir /path/to/notifico
exec ./env/bin/python -m notifico www 

